# supervisor-conf.d/bsp.conf
[program:bsp]
command = bsp --start
logfile=/var/log/bsp.log
autostart=true
autorestart=true
startsecs=0