# supervisor-conf.d/bsp_server.conf
[program:bsp_server]
command = bsp_server
logfile=/var/log/bsp_server.log
autostart=true
autorestart=true
startsecs=3